# deployment/com.lingxing.etl.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.lingxing.etl</string>

    <key>WorkingDirectory</key>
    <string>/Users/wangjianhuang/Projects/lingxing_erp_exporter</string>

    <key>Program</key>
    <string>/Users/wangjianhuang/Projects/lingxing_erp_exporter/run_daily.sh</string>

    <key>StandardOutPath</key>
    <string>/Users/wangjianhuang/Library/Logs/lingxing-etl.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/wangjianhuang/Library/Logs/lingxing-etl-error.log</string>

    <!-- 每小时触发一次，由 run_daily.sh 内部判断是否到执行时间 -->
    <!-- 这样无需手动切换夏令时/冬令时 -->
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Hour</key>
            <integer>15</integer>
            <key>Minute</key>
            <integer>20</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>16</integer>
            <key>Minute</key>
            <integer>20</integer>
        </dict>
    </array>

    <!-- 开机/登录时立即触发一次，由 run_daily.sh 内部判断是否需要执行 -->
    <key>RunAtLoad</key>
    <true/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>OPENCLAW_HOOK_TOKEN</key>
        <string>8394dfcc44c1c7bc58a379a348ab2f8ef4ecc3a6c9938c1f</string>
        <key>OPENCLAW_PATROL_GROUP</key>
        <string>oc_c142005c008b351c8587f8f8659877b8</string>
        <key>OPENCLAW_REPORT_GROUP</key>
        <string>oc_77700e43d8bc2fc7bef35272f30048e8</string>
    </dict>

    <key>KeepAlive</key>
    <false/>
</dict>
</plist>
