# deployment/com.lingxing.etl.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.lingxing.etl</string>

    <key>WorkingDirectory</key>
    <string>/Users/wangjianhuang/Projects/lingxing_erp_exporter</string>

    <key>Program</key>
    <string>/Users/wangjianhuang/Projects/lingxing_erp_exporter/run_daily.sh</string>

    <key>StandardOutPath</key>
    <string>/Users/wangjianhuang/Library/Logs/lingxing-etl.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/wangjianhuang/Library/Logs/lingxing-etl-error.log</string>

    <!-- 每小时触发一次，由 run_daily.sh 内部判断是否到执行时间 -->
    <!-- 这样无需手动切换夏令时/冬令时 -->
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Hour</key>
            <integer>15</integer>
            <key>Minute</key>
            <integer>30</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>16</integer>
            <key>Minute</key>
            <integer>30</integer>
        </dict>
    </array>

    <key>RunAtLoad</key>
    <false/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>OPENCLAW_HOOK_TOKEN</key>
        <string>YOUR_HOOK_TOKEN</string>
        <key>OPENCLAW_PATROL_GROUP</key>
        <string>oc_PATROL_GROUP_ID</string>
        <key>OPENCLAW_REPORT_GROUP</key>
        <string>oc_REPORT_GROUP_ID</string>
    </dict>

    <key>KeepAlive</key>
    <false/>
</dict>
</plist>
